# xianyu/agent_mvc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# Tornado应用进程组
[group:tornado]
programs=tornado-8000,tornado-8001,tornado-8002,tornado-8003

# Tornado实例1
[program:tornado-8000]
command=python app.py --port=8000
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado-8000.log
stderr_logfile=/var/log/tornado-8000-error.log
environment=PYTHONPATH="/app",PORT="8000"

# Tornado实例2
[program:tornado-8001]
command=python app.py --port=8001
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado-8001.log
stderr_logfile=/var/log/tornado-8001-error.log
environment=PYTHONPATH="/app",PORT="8001"

# Tornado实例3
[program:tornado-8002]
command=python app.py --port=8002
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado-8002.log
stderr_logfile=/var/log/tornado-8002-error.log
environment=PYTHONPATH="/app",PORT="8002"

# Tornado实例4
[program:tornado-8003]
command=python app.py --port=8003
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado-8003.log
stderr_logfile=/var/log/tornado-8003-error.log
environment=PYTHONPATH="/app",PORT="8003"
